# systemd/anki-conjugation.timer
[Unit]
Description=Run Anki conjugation unlocker daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
